# scripts/init/easylab-server.conf
start on startup
task
exec /home/smith/startEasylabServer.sh
